# systemd/ampelmann.service
[Unit]
Description=Ampelmann system monitor
After=network.target ollama.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ampelmann run
ExecStartPost=/usr/local/bin/ampelmann dashboard
User=root
Nice=10
IOSchedulingClass=idle
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ampelmann

# Security hardening
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/lib/ampelmann /var/log/ampelmann /srv/site/ampelmann/www/data

[Install]
WantedBy=multi-user.target
